# libreoffice/src/deb/control/prerm
#!/bin/sh

if [ -x /bin/systemctl ]; then
   systemctl stop soffice.service
   systemctl disable soffice.service
fi